# deploy/systemd/face-service.service
[Unit]
Description=MyFirstProject Face Recognition Service (Python Flask)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/myapp

# Ensure opencv-python-headless is used (no GUI/display available on servers)
# This mirrors the Dockerfile fix and survives container/service restarts
ExecStartPre=/var/www/myapp/venv/bin/pip uninstall opencv-python -y
ExecStartPre=/var/www/myapp/venv/bin/pip install opencv-python-headless --quiet

# Use the virtualenv Python
ExecStart=/var/www/myapp/venv/bin/python face_service.py

Environment="FACE_SERVICE_PORT=5050"

Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=face-service

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
